# rootfs/etc/services.d/tdarr-node/finish
#!/usr/bin/execlineb -P
echo "Tdarr Node stopped."
