# rootfs/etc/services.d/tdarr-node/run
#!/usr/bin/execlineb -P
with-contenv
cd /opt/tdarr/Tdarr_Node
node Tdarr_Node.js
